# crontabs
0 1 * * * bash /root/pr-repo/update-repo-packages.sh